# freecad/cleanup.cron
0 * * * * find /var/www/html/out_stl/*.stl -mmin +1440 -delete 

